# files/system/maetel/usr/etc/profile.d/rust_alias.sh
#/bin/bash
eval "$(zoxide init bash)" 
alias cd=z
alias find=fd
alias vim=helix
alias grep=rg